# check_mk_agent/files/add-to-cmk.conf
description     "add host to check_mk"

start on (local-filesystems and net-device-up IFACE!=lo)

task

exec /usr/local/bin/add-host-check_mk.sh

